# kooiot/factory_tlink/files/etc/rc.local

# Timeout time ( seconds )
TIMEOUT=3

echo "Flashing eMMC will be started after ${TIMEOUT}!"
echo "Enter any character followed by <ENTER> to stop flashing"
read -t ${TIMEOUT} input_stop

if [ -z "${input_stop}" ];
then
	/sbin/factory_make &
else
	echo "Flashing stoped!!!!"
fi

exit 0
